# discord-bot.service
[Unit]
Description=Panda Discord LLM Bot
Documentation=https://github.com/jcpelletier/panda
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=discord-bot
Group=discord-bot
WorkingDirectory=/opt/discord-bot
EnvironmentFile=/opt/discord-bot/.env
ExecStart=/opt/discord-bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot

# Hardening — the bot needs network, /proc, Docker socket, log files, and media library
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log /opt/discord-bot /mnt/media

# Ambient capabilities let smartctl subprocesses read raw SMART data.
# NoNewPrivileges is intentionally omitted — it blocks ambient cap inheritance.
AmbientCapabilities=CAP_SYS_RAWIO CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
